# scripts/com.buckmountain.farm.media-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.buckmountain.farm.media-watcher</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/iamclaw/Library/Application Support/buckmountain-farm/openclaw_watcher.py</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>EnvironmentVariables</key>
    <dict>
        <key>BUCKMOUNTAIN_DASHBOARD_URL</key>
        <string>https://buckmountain-farm.vercel.app/api/admin/assets</string>
        <key>BUCKMOUNTAIN_DASHBOARD_TOKEN</key>
        <string>2574cf2e3433917aebee5938488d7dfca7b001b4eafd39c141dd637a759a829c</string>
        <key>BUCKMOUNTAIN_VERCEL_BYPASS</key>
        <string>X2fz2zpuN9fZ2xyJU6zDRSBQUSCFIyjR</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/iamclaw/Library/Logs/buckmountain-farm-watcher.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/iamclaw/Library/Logs/buckmountain-farm-watcher.err.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/iamclaw/Library/Application Support/buckmountain-farm</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
